# conf/mctp-local.target
[Unit]
Description=MCTP local network configuration
